# s6-overlay/s6-rc.d/gunicorn/run
#!/command/execlineb -P
cd ~/ptk-admin-panel
uv run gunicorn -w 4 -b 0.0.0.0:5000 src.ptk_admin_panel.app:app